# deploy/com.e-newspaper.daily-update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Daily newspaper update job.

        Retry ladder at 14:30 / 15:30 / 16:30 / 17:30 HKT every day.
        The NYT next-day print PDF uploads overnight ET, which lands
        ~11:00–15:00 HKT same day — occasionally later on weekends. Four
        attempts an hour apart catches the edition as soon as it appears;
        nyt_scraper.py is idempotent (writes a the-new-york-times.etdate
        marker after a successful fetch) so the later runs become no-ops
        once the day's edition has been captured.

        The pipeline each run:
          1. scraper.py downloads today's papers from frontpages.com
             (and high-res sources for papers with dedicated scrapers)
          2. processor.py processes the chosen paper for the e-ink display

        The OpenDisplay server (separate launchd job) picks up the new
        images/current.png on its next client poll.

        Install:
          cp deploy/com.e-newspaper.daily-update.plist ~/Library/LaunchAgents/
          launchctl load ~/Library/LaunchAgents/com.e-newspaper.daily-update.plist

        Manually trigger a run:
          launchctl start com.e-newspaper.daily-update

        Tail logs:
          tail -f <install-dir>/scraper.log

        Note: deploy/install.sh rewrites __INSTALL_DIR__ to match the
        actual install location before copying this file into ~/Library/LaunchAgents.
    -->
    <key>Label</key>
    <string>com.e-newspaper.daily-update</string>

    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <!--
        launchd runs with a stripped-down PATH that does NOT include
        /opt/homebrew/bin, so Homebrew-installed tools (pdftoppm from
        poppler, in particular) are invisible. Without this block,
        nyt_scraper silently falls back to the 600x800 frontpages.com
        thumbnail. Don't remove this unless you also make the Python
        code invoke pdftoppm by absolute path.
    -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/sh</string>
        <string>-c</string>
        <string>.venv/bin/python src/scraper.py &gt;&gt; scraper.log 2&gt;&amp;1 &amp;&amp; .venv/bin/python src/processor.py &gt;&gt; scraper.log 2&gt;&amp;1</string>
    </array>

    <!--
        Fire every hour (3600 seconds). StartInterval is timezone-
        independent, unlike StartCalendarInterval which on this Mac
        mini has been observed caching PDT despite the system being
        set to Asia/Hong_Kong — the original 14:30-17:30 HKT schedule
        was actually firing at 05:30-08:30 HKT (14:30-17:30 PDT offset
        by 9 hours). Rather than hack the hours to account for the
        cached TZ, just fire hourly: the scraper is idempotent (NYT
        via the etdate marker, frontpages.com is cheap), so the bulk
        of runs are ~3s no-ops. Cost is minor log spam; benefit is
        catching a fresh NYT within an hour of upload regardless of
        what TZ launchd thinks it's in.
    -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/launchd-update.out.log</string>
    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/launchd-update.err.log</string>
</dict>
</plist>
