# deploy/com.e-newspaper.auto-recover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Auto-recovery daemon for the e-ink display (heartbeat-triggered).

        Listens on TCP 9999 for the display's reverse-shell dial-in. The
        server heartbeat (images/last-poll.txt) is the only reliable failure
        signal: during the ENETUNREACH bug the display's own dumpsys reports
        a healthy network, so the display cannot self-detect it. When the
        heartbeat goes stale this daemon sends ONE eth0-bounce + app-restart
        over the reverse shell, then waits 15 min before acting again, and
        backs off to hourly after 3 attempts. That rate limit is the whole
        point: bouncing eth0 every ~5 min is what wedged the network and
        caused the 2026-06-06 outage.

        Install:
          cp deploy/com.e-newspaper.auto-recover.plist ~/Library/LaunchAgents/
          launchctl load ~/Library/LaunchAgents/com.e-newspaper.auto-recover.plist

        Tail recovery events / heartbeat state:
          tail -F auto-recover.err.log | grep -E "STALE|recovery|healthy"

        Stop temporarily (e.g. to use tools/remote_shell.py manually):
          launchctl unload ~/Library/LaunchAgents/com.e-newspaper.auto-recover.plist

        Note: deploy/install.sh rewrites __INSTALL_DIR__ before copying.
    -->
    <key>Label</key>
    <string>com.e-newspaper.auto-recover</string>

    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__INSTALL_DIR__/.venv/bin/python</string>
        <string>-u</string>
        <string>src/auto_recover.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/auto-recover.log</string>
    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/auto-recover.err.log</string>
</dict>
</plist>
