# deploy/com.e-newspaper.auto-recover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        Auto-recovery daemon for the e-ink display (heartbeat-triggered).

        Listens on TCP 9999 for the display's reverse-shell dial-in. The
        server heartbeat (images/last-poll.txt) is the only reliable failure
        signal: during the failure the display's own dumpsys can report a
        healthy network, so the display cannot self-detect it. When the
        heartbeat goes stale this daemon REBOOTS the display over the reverse
        shell, waits 10 min, and backs off to hourly after 4 reboots. A plain
        reboot clears both observed failure forms (app ENETUNREACH, and no
        network agent at all); earlier eth0-bounce versions fixed only one,
        caused the other, and fought tp_watchdog.

        Install:
          cp deploy/com.e-newspaper.auto-recover.plist ~/Library/LaunchAgents/
          launchctl load ~/Library/LaunchAgents/com.e-newspaper.auto-recover.plist

        Tail recovery events / heartbeat state:
          tail -F auto-recover.err.log | grep -E "STALE|recovery|healthy"

        Stop temporarily (e.g. to use tools/remote_shell.py manually):
          launchctl unload ~/Library/LaunchAgents/com.e-newspaper.auto-recover.plist

        Note: deploy/install.sh rewrites __INSTALL_DIR__ before copying.
    -->
    <key>Label</key>
    <string>com.e-newspaper.auto-recover</string>

    <key>WorkingDirectory</key>
    <string>__INSTALL_DIR__</string>

    <key>ProgramArguments</key>
    <array>
        <string>__INSTALL_DIR__/.venv/bin/python</string>
        <string>-u</string>
        <string>src/auto_recover.py</string>
    </array>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>__INSTALL_DIR__/auto-recover.log</string>
    <key>StandardErrorPath</key>
    <string>__INSTALL_DIR__/auto-recover.err.log</string>
</dict>
</plist>
